# scripts/kode-kronical-daemon.service
[Unit]
Description=kode-kronical System Monitoring Daemon
Documentation=https://github.com/jeremycharlesgillespie/kode-kronical
After=network.target

[Service]
Type=simple
# Use the venv's kode-kronical-daemon directly (it uses the venv's Python internally)
# Replace __USER__ with actual username during installation
ExecStart=/home/__USER__/.venv/bin/kode-kronical-daemon -c /etc/kode-kronical/daemon.yaml start --foreground
ExecStop=/home/__USER__/.venv/bin/kode-kronical-daemon -c /etc/kode-kronical/daemon.yaml stop
ExecReload=/home/__USER__/.venv/bin/kode-kronical-daemon -c /etc/kode-kronical/daemon.yaml restart
KillMode=process
Restart=on-failure
RestartSec=10
User=__USER__
Environment=HOME=/home/__USER__
WorkingDirectory=/home/__USER__

# Security settings (relaxed for compatibility)
ReadWritePaths=/home/__USER__/.local /home/__USER__/.config /etc/kode-kronical

# Resource limits
LimitNOFILE=4096
CPUQuota=20%
MemoryMax=256M

[Install]
WantedBy=multi-user.target